# board/canaan/k510/k510_rootfs_skeleton/etc/inittab
# system initialisation
::sysinit:/etc/init.d/rc.sysinit

# run gettys on the serial ports
::respawn:-/sbin/getty -n -L ttyS0 115200 vt100 -l usr/autologin.profile
ttyS0::sysinit:etc/init.d/run.sh

# stuff to do before rebooting
::ctrlaltdel:/etc/init.d/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
